# root/etc/rc.d/rc.local
#!/bin/bash
sleep 20; /usr/bin/xs-vpn test
